# scripts/clickme.service.example
[Unit]
Description=Clickme Next.js application
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=clickme
Group=clickme
WorkingDirectory=/srv/clickme/current
EnvironmentFile=/srv/clickme/shared/clickme.env
Environment=NEXT_TELEMETRY_DISABLED=1
# Nginx owns 127.0.0.1:3000. The explicit env command is the final authority
# for the private Next.js listener even if the shared env file still has an old
# PORT or HOSTNAME value.
ExecStart=/usr/bin/env NODE_ENV=production PORT=3001 HOSTNAME=127.0.0.1 /opt/clickme-node-v24.18.0/bin/node /srv/clickme/current/server.js
Restart=always
RestartSec=5s
TimeoutStopSec=20s
KillSignal=SIGTERM
UMask=0077
MemoryHigh=768M
MemoryMax=1G
TasksMax=128
LimitNOFILE=8192
OOMPolicy=stop
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clickme
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
CapabilityBoundingSet=
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
