# scripts/clickme-cleanup.timer.example
[Unit]
Description=Run Clickme operational-data cleanup hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true
Unit=clickme-cleanup.service

[Install]
WantedBy=timers.target
